# systemd/user/gpg-agent.service.d/override.conf
[Service]
ExecStartPost=/usr/bin/echo UPDATESTARTUPTTY | /usr/bin/gpg-connect-agent