# systemd/user/redshift.service.d/override.conf
[Service]
Environment=DISPLAY=:0.0
Restart=
Restart=no
