# systemd/user/podman.service.d/override.conf
[Service]
Environment=LOGGING="--log-level=warning"